# resources/etc/systemd/system/host-setup.service
[Unit]
Description=Host system setup
After=systemd-machined.service

[Service]
Type=oneshot
RemainAfterExit=yes
ExecStart=/sbin/host-setup start
ExecStop=/sbin/host-setup stop

[Install]
WantedBy=multi-user.target
